# rango-gate.timer
[Unit]
Description=Nightly LaMP-2 gate pipeline (timer)

[Timer]
OnCalendar=*-*-* 03:30
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
